# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:redis]
command=redis-server --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru --requirepass %(ENV_REDIS_PASSWORD)s --bind 127.0.0.1 --port 6379
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/opt/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app
user=appuser
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/app",
    STUDIO_DATABASE_URL="sqlite+aiosqlite:///data/studio.db",
    STUDIO_REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    STUDIO_AGENT_SERVICE_TOKEN="%(ENV_STUDIO_AGENT_SERVICE_TOKEN)s",
    STUDIO_GENERATED_DIR="/shared/generated",
    STUDIO_ENCRYPTION_KEY="%(ENV_STUDIO_ENCRYPTION_KEY)s",
    STUDIO_JWT_SECRET="%(ENV_STUDIO_JWT_SECRET)s"

[program:mcp]
command=/app/mcp_entrypoint.sh
directory=/app
user=appuser
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/shared/generated",
    REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    GENERATED_DIR="/shared/generated"

[program:agent-mcp]
command=/opt/venv/bin/python -m agent_mcp
directory=/app
user=appuser
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/app",
    STUDIO_BACKEND_URL="http://127.0.0.1:8000",
    STUDIO_REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    STUDIO_AGENT_SERVICE_TOKEN="%(ENV_STUDIO_AGENT_SERVICE_TOKEN)s"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
