# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:redis]
command=redis-server --dir /data/redis --appendonly yes --save 300 10 --save 60 1000 --maxmemory 256mb --maxmemory-policy volatile-lru --requirepass %(ENV_REDIS_PASSWORD)s --bind 127.0.0.1 --port 6379
user=redis
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/opt/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app
user=appuser
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/app",
    STUDIO_DATABASE_URL="sqlite+aiosqlite:///data/studio.db",
    STUDIO_REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    STUDIO_AGENT_SERVICE_TOKEN="%(ENV_STUDIO_AGENT_SERVICE_TOKEN)s",
    STUDIO_GENERATED_DIR="/shared/generated",
    STUDIO_ENCRYPTION_KEY="%(ENV_STUDIO_ENCRYPTION_KEY)s",
    STUDIO_JWT_SECRET="%(ENV_STUDIO_JWT_SECRET)s",
    STUDIO_NGINX_CONFIG_DIR="/etc/nginx/conf.d",
    STUDIO_CERTBOT_WEBROOT="/var/www/certbot",
    STUDIO_LETSENCRYPT_DIR="/etc/letsencrypt",
    STUDIO_SSL_STAGING="%(ENV_STUDIO_SSL_STAGING)s",
    STUDIO_EXTERNAL_HTTPS_PORT="%(ENV_STUDIO_EXTERNAL_HTTPS_PORT)s",
    STUDIO_NGINX_BACKEND_HOST="127.0.0.1",
    STUDIO_NGINX_MCP_HOST="127.0.0.1",
    STUDIO_NGINX_AGENT_MCP_HOST="127.0.0.1",
    STUDIO_NGINX_HTTP2_MODERN="false"

[program:mcp]
command=/app/mcp_entrypoint.sh
directory=/app
user=appuser
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/shared/generated",
    REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    GENERATED_DIR="/shared/generated"

[program:agent-mcp]
command=/opt/venv/bin/python -m agent_mcp
directory=/app
user=appuser
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    HOME="/home/appuser",
    PYTHONPATH="/app",
    STUDIO_BACKEND_URL="http://127.0.0.1:8000",
    STUDIO_REDIS_URL="redis://:%(ENV_REDIS_PASSWORD)s@127.0.0.1:6379/0",
    STUDIO_AGENT_SERVICE_TOKEN="%(ENV_STUDIO_AGENT_SERVICE_TOKEN)s"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx-reload-watcher]
command=/bin/sh -c 'while true; do if [ -f /etc/nginx/conf.d/.reload ]; then if nginx -t 2>/dev/null; then nginx -s reload && rm -f /etc/nginx/conf.d/.reload; else echo "nginx config test failed, skipping reload" >&2; rm -f /etc/nginx/conf.d/.reload; fi; fi; sleep 2; done'
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
